# scripts/com.health-logger.dev-env-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS launchd による定期実行設定（月1回）
  
  セットアップ手順:
  1. このファイルを編集して <YOUR_HOME_DIR> と <REPO_PATH> を実際のパスに置換
  2. ~/Library/LaunchAgents/ にコピー
     cp scripts/com.health-logger.dev-env-review.plist \
        ~/Library/LaunchAgents/com.health-logger.dev-env-review.plist
  3. launchd に登録
     launchctl load ~/Library/LaunchAgents/com.health-logger.dev-env-review.plist
  4. 確認
     launchctl list | grep health-logger
  5. 手動実行テスト
     launchctl start com.health-logger.dev-env-review
  
  停止:
     launchctl unload ~/Library/LaunchAgents/com.health-logger.dev-env-review.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.health-logger.dev-env-review</string>

  <!-- 実行コマンド: bash scripts/run-dev-env-review.sh -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string><REPO_PATH>/scripts/run-dev-env-review.sh</string>
    <string>--model</string>
    <string>haiku</string>
  </array>

  <!-- 作業ディレクトリ -->
  <key>WorkingDirectory</key>
  <string><REPO_PATH></string>

  <!-- 実行スケジュール: 毎月1日 AM 9:00 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Day</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- ログ出力先 -->
  <key>StandardOutPath</key>
  <string><YOUR_HOME_DIR>/.claude/logs/dev-env-review.log</string>
  <key>StandardErrorPath</key>
  <string><YOUR_HOME_DIR>/.claude/logs/dev-env-review.err</string>

  <!-- ログイン時に自動起動 -->
  <key>RunAtLoad</key>
  <false/>

  <!-- 失敗時の再起動なし（誤課金防止）-->
  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
